# v2/systemd/autonome-review.service
[Unit]
Description=Autonome v2 LLM Review Report

[Service]
Type=oneshot
WorkingDirectory=/mnt/e/NomadCrew[GROWTH]/trading-os/v2
ExecStart=/bin/bash -lc 'cd /mnt/e/NomadCrew[GROWTH]/trading-os/v2 && source .venv/bin/activate && python3 -m autonome.supervisor.review > /mnt/e/NomadCrew[GROWTH]/trading-os/v2/data/review_$(date +%%Y%%m%%d_%%H%%M).txt'
